# etc/systemd/camera-node.service
[Unit]
Description=Rover Camera ROS2 Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rover
Group=rover
WorkingDirectory=/home/rover/camera_ws
Environment=ROS_DISTRO=lyrical
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_cpp
Environment=ROS_DOMAIN_ID=0
Environment=LD_LIBRARY_PATH=/home/rover/camera_ws/install/libcamera/lib:/home/rover/camera_ws/install/libcamera/lib/aarch64-linux-gnu
ExecStart=/bin/bash -lc 'source /opt/ros/lyrical/setup.bash && source /home/rover/camera_ws/install/setup.bash && exec ros2 run camera_ros camera_node --ros-args --params-file /home/rover/.ros/camera_node.yaml'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
